# scripts/launchd/curaitor.langgraph.sample.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.curaitor.langgraph.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>python</string>
    <string>scripts/run_daily.py</string>
    <string>--query</string>
    <string>plastic recycling</string>
    <string>--max-days</string>
    <string>7</string>
    <string>--db</string>
    <string>/Users/riteshkumar/Library/CloudStorage/Box-Box/Research-postdoc/oxRSE-project/oxRSE-2025/data-codes-sharing/curaitor-agent/data/curaitor.sqlite</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/riteshkumar/Library/CloudStorage/Box-Box/Research-postdoc/oxRSE-project/oxRSE-2025/data-codes-sharing/curaitor-agent</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/riteshkumar/Library/CloudStorage/Box-Box/Research-postdoc/oxRSE-project/oxRSE-2025/data-codes-sharing/logs/langgraph_daily.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/riteshkumar/Library/CloudStorage/Box-Box/Research-postdoc/oxRSE-project/oxRSE-2025/data-codes-sharing/logs/langgraph_daily.err</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>

